# build/boardingpass-bt@.service
[Unit]
Description=BoardingPass Bluetooth PAN on %i
PartOf=boardingpass.service

[Service]
Type=simple
ExecStartPre=/usr/bin/bluetoothctl power on
ExecStartPre=/usr/bin/bluetoothctl discoverable on
ExecStart=/usr/lib/boardingpass/scripts/bt-pan.sh %i br-bp0
ExecStartPost=/usr/sbin/ip addr add 10.0.1.1/24 dev br-bp0
ExecStopPost=/usr/bin/bluetoothctl discoverable off
Restart=no
